# setupfiles/adminserver.service
# AdminServer systemctl service file
#
# Link and enable with:
# % sudo ln -sf /home/debian/driftconditions/setupfiles/adminserver.service /etc/systemd/system/adminserver.service
# % sudo systemctl daemon-reload
# % sudo systemctl enable adminserver.service   # 'linked' won't auto-start on boot; enable creates the wants symlink
#
# Start:
# % sudo systemctl start adminserver.service
#
# Restart:
# % sudo systemctl restart adminserver.service
#
# View journal:
# % journalctl -u adminserver.service --no-pager | less
# (though ordinary server messages should be reported in driftconditions/logs/adminserver.log)

[Unit]
Description=AdminServer Node.js Express backend service
# mysql.service uses Type=notify so systemd waits for MySQL to be fully ready before starting this
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=debian
WorkingDirectory=/home/debian/driftconditions/AdminServer
ExecStart=/usr/bin/npm start
Restart=on-failure

[Install]
WantedBy=multi-user.target
