# setupfiles/audio-analysis.service
# Audio analysis runner systemd service file
# Invoked by audio-analysis.timer — do not start directly.
#
# Link into place with:
# % sudo ln -sf /home/debian/driftconditions/setupfiles/audio-analysis.service /etc/systemd/system/audio-analysis.service
# % sudo ln -sf /home/debian/driftconditions/setupfiles/audio-analysis.timer /etc/systemd/system/audio-analysis.timer
# % systemctl daemon-reload
# % sudo systemctl enable --now audio-analysis.timer
#
# Check timer status:
# % systemctl status audio-analysis.timer
#
# View logs:
# % journalctl -u audio-analysis.service --no-pager | less

[Unit]
Description=DriftConditions audio analysis runner
After=network.target

[Service]
Type=oneshot
User=debian
WorkingDirectory=/home/debian/driftconditions/AdminServer
Environment=NODE_PATH=/home/debian/driftconditions/AdminServer/node_modules
ExecStart=/usr/bin/node /home/debian/driftconditions/scripts/run-audio-analysis.js

[Install]
WantedBy=multi-user.target
