# setupfiles/liquidsoap.service
[Unit]
Description=Liquidsoap Streaming Service
After=network.target mixengine.service
Wants=mixengine.service

[Service]
User=liquidsoap
Group=liquidsoap
ExecStart=/bin/liquidsoap /home/debian/driftconditions/setupfiles/liquidsoap.liq
Restart=on-failure

[Install]
WantedBy=multi-user.target